# shopping/custom_configs/elasticsearch.ini
[program:elasticsearch]
command=/usr/share/java/elasticsearch/bin/elasticsearch
environment=ES_JAVA_HOME=/usr,ES_JAVA_OPTS="-Xms512m -Xmx512m"
autostart=true
autorestart=true
priority=8
startretries=3
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
