# shopping/custom_configs/mysql.ini
# for supervisord — overrides the SIF's /etc/supervisor.d/mysql.ini
# Redirects mysqld output to a real log file so crashes can be diagnosed.
[program:mysqld]
command=mysqld --log-error=/var/log/mysql/error.log
autostart=true
autorestart=true
priority=1
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/mysql/mysqld.out
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile=/var/log/mysql/mysqld.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=1
